# deploy/pizza-celery-beat.service
[Unit]
Description=Pizza Scheduler Celery Beat (Scheduler)
After=network.target redis.service

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/pizza
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=/home/ubuntu/pizza/.env.production
ExecStart=/home/ubuntu/pizza/venv/bin/celery -A config beat -l INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/celery/beat.log
StandardError=append:/var/log/celery/beat-error.log

[Install]
WantedBy=multi-user.target
